# roles/network-wifi-fallback/files/wifi-ap-fallback.service
[Unit]
Description=Wi-Fi Access Point Fallback Monitor
After=network.target systemd-networkd.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi-ap-fallback.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
