# roles/network-fallback/files/network-fallback.service
[Unit]
Description=Network Fallback DHCP Monitor
After=network.target systemd-networkd.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/network-fallback.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
